# deb/etc/init/umountfstab.conf
start on runlevel [06]

task

script
	make -f /etc/mountfstab/Makefile umount
end script
