# deb/etc/init/umountfstab.conf
start on runlevel [06]

task

script
	make -C /etc/mountfstab/ umount 2>&1 | logger -t umountfstab
end script
